# conf/sample/supervisord_sample.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[unix_http_server]
file={{ PROJECT_DIR }}/var/run/supervisor.sock   ; (the path to the socket file)


[supervisord]
logfile={{ PROJECT_DIR }}/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile={{ PROJECT_DIR }}/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)


; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[supervisorctl]
serverurl=unix://{{ PROJECT_DIR }}/var/run/supervisor.sock ; use a unix:// URL  for a unix socket


[program:unisample]
stdout_logfile={{ PROJECT_DIR }}/var/log/projectname_{{ PROFILE }}_site.log
stderr_logfile={{ PROJECT_DIR }}/var/log/projectname_{{ PROFILE }}_site_err.log
command={{ VIRTUALENV_BIN }}/gunicorn -c {{ PROJECT_DIR }}/conf/gunicorn_conf.py conf.wsgi:application
directory={{ PROJECT_DIR }}
user=www-data
autostart=true
autorestart=true
